# unixdomain/srv.socket
[Socket]
ListenStream=/run/srv.sk

[Install]
WantedBy=sockets.target
