# unixdomain/srv.service
[Service]
ExecStart=/usr/bin/srvd

[Install]
WantedBy=multi-user.target
